# systemd/vllm-minimax.service
[Unit]
Description=vLLM serving MiniMax-M2.5 (8xH20)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=/etc/llm-deploy.env
WorkingDirectory=/AI4S/Users/howardwang/llm-deploy
ExecStart=/AI4S/Users/howardwang/llm-deploy/scripts/start_vllm.sh
Restart=on-failure
RestartSec=15
# 模型加载 + CUDA graph 编译可能需要几分钟
TimeoutStartSec=1800
TimeoutStopSec=120
KillSignal=SIGINT
LimitNOFILE=1048576
# vLLM 会 fork 多个 worker 进程，放开 nproc
LimitNPROC=infinity
# 让 systemd 把整棵进程树一起收掉
KillMode=mixed

[Install]
WantedBy=multi-user.target
